# deploy/rag-server.svc.plus/systemd/rag-server.service
[Unit]
Description=rag-server API
After=network-online.target stunnel-pg.service
Wants=network-online.target stunnel-pg.service

[Service]
Type=simple
Environment=DATABASE_URL=postgres://user:pass@127.0.0.1:15432/rag?sslmode=disable
ExecStart=/usr/local/bin/rag-server --config /etc/rag-server/server.yaml
Restart=on-failure
RestartSec=3s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
